# supervisord-deploy.conf
[program:deploy]
command=/deploy.sh dblol root toor secretlol
numprocs=1
autostart=true
autorestart=false
